# tachiyomi-backup-tracker.path
[Unit]
Description=Watch Tachiyomi backup

[Path]
PathModified=/path/to/backup
Unit=tachiyomi-backup-tracker.service

[Install]
WantedBy=multi-user.target